# ssmtp.service.in
[Unit]
Description=sSMTP-Reloaded timer service
Documentation=man:ssmtp
RefuseManualStart=yes
RefuseManualStop=yes

[Service]
Type=oneshot
ExecStart=SSMTP -q
