# ssmtp.timer
[Unit]
Description=sSMTP-Reloaded Timer
Documentation=man:ssmtp

[Timer]
Persistent=true
OnUnitActiveSec=5min
Unit=ssmtp.service
