# ssmtp.service.in
[Unit]
Description=sSMTP-Reloaded timer service
Documentation=man:ssmtp

[Service]
Type=oneshot
ExecStart=SSMTP
